# jibri/rootfs/etc/services.d/25-pulseaudio/run
#!/usr/bin/with-contenv bash
HOME=/home/jibri

DAEMON="/usr/bin/start-pulseaudio-x11"
exec s6-setuidgid jibri /bin/bash -c "exec $DAEMON"
